# modules/prometheus/templates/initscripts/prometheus-node-exporter.upstart.erb
description "prometheus-node-exporter"
respawn
respawn limit 10 30

start on (local-filesystem and net-device-up IFACE!=lo)

setuid prometheus

script
    . /etc/default/prometheus-node-exporter
    exec /usr/bin/prometheus-node-exporter $ARGS
end script
