# modules/wikimetrics/templates/upstart.wikimetricsqueue.conf.erb
start on startup
script
    cd <%= @wikimetrics_path %>
    wikimetrics --mode queue --db-config /etc/wikimetrics/db_config.yaml --queue-config /etc/wikimetrics/queue_config.yaml
end script
