# modules/openstack/files/antelope/designate/designate-worker.conf
# SPDX-License-Identifier: Apache-2.0

# This init scripts are not present in the Newton debian packages.
#  Probably we won't need this on Ocata.

description "Designate worker"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

setuid designate

chdir /var/lib/designate
exec /usr/bin/designate-worker --config-file=/etc/designate/designate.conf --log-file=/var/log/designate/designate-worker.log
respawn
respawn limit 10 30
