# modules/keyholder/files/keyholder-agent.conf
# keyholder-agent - Shared SSH-agent
#
# Runs the ssh-agent(1) instance that holds shared identities.

description "Shared SSH agent"

start on (local-filesystems and net-device-up IFACE!=lo)

setgid keyholder
setuid keyholder

exec /usr/bin/ssh-agent -d -a /run/keyholder/agent.sock
post-start exec [ -S /run/keyholder/agent.sock ] || sleep 1
post-stop exec /bin/rm -f /run/keyholder/agent.sock

# vim: set ft=upstart:
