# modules/keyholder/files/keyholder-proxy.conf
# keyholder-proxy - Filtering proxy for ssh-agent(1)
#
# The `keyholder-proxy` service runs the filtering ssh-agent proxy
# that acts as an intermediary between users in the trusted group
# and the backend ssh-agent that holds the shared key(s).

description "Shared SSH agent proxy"

start on started keyholder-agent
stop on stopped keyholder-agent

pre-start script
    [ ! -S /run/keyholder/agent.sock ] && sleep 0.5
    [ ! -r /etc/default/keyholder ] && { stop; exit 0; }
    . /etc/default/keyholder
    /bin/rm -f /run/keyholder/proxy.sock
end script

script
  . /etc/default/keyholder
  /sbin/start-stop-daemon --quiet --start    \
    --chuid "keyholder:${KEYHOLDER_GROUP}"   \
    --pidfile /run/keyholder/proxy.pid       \
    --umask 007                              \
    --startas /usr/local/bin/ssh-agent-proxy
end script

post-stop exec /bin/rm -f /run/keyholder/proxy.sock

# vim: set ft=upstart:
