# modules/mwprof/templates/upstart/profiler-to-carbon.conf.erb
# mwprof
#
# MediaWiki profiling data collector
#
description "profiling data collector"

start on started mwprof/collector
stop on runlevel [!2345] or mwprof.stop

setuid mwprof
setgid mwprof

exec /srv/deployment/mwprof/mwprof/profiler-to-carbon \
    --collector-host=127.0.0.1 \
    --collector-port=<%= @collector_port %> \
    --carbon-host=<%= @carbon_host %> \
    --carbon-port=<%= @carbon_port %>
