# modules/interface/templates/enable-rps.conf.erb
# enable-rps

description "Enable RPS on <%= @interface %> RX queues"

start on filesystem
task
script
	interface-rps <%= @interface %>
end script
