# modules/interface/templates/enable-rps.conf.erb
# enable-rps

description "Enable RPS/RSS on <%= @interface %> RX queues"

start on filesystem
task
script
	interface-rps <%= @interface %> <%= @rss_pattern %>
end script
