# modules/toollabs/files/iptables.conf
# iptables - load /etc/iptables.conf if it is present

description "Load /etc/iptables.conf"

start on net-device-up

respawn

console log

pre-start script
    if [ -r /etc/iptables.conf ]; then
        /sbin/iptables-restore /etc/iptables.conf
    fi
end script

