# modules/toollabs/files/toolwatcher.conf
description "creates tool homes"

start on runlevel [2345]
stop on runlevel [06]

respawn

exec /usr/sbin/toolwatcher
