# modules/labstore/templates/initscripts/cleanup.timer.erb
[Unit]
Description=Cleanup older snapshots

[Timer]
OnCalendar=01:00
Persistent=true

[Install]
WantedBy=timers.target
