# modules/labstore/templates/initscripts/cleanup.timer.erb
[Unit]
Description=Cleanup snapshots (<%= @title %>)

[Timer]
OnCalendar=<%= @calendar %>
Persistent=true

[Install]
WantedBy=timers.target
