# modules/labstore/templates/initscripts/replicate.timer.erb
[Unit]
Description=Replicate storage between hosts (<%= @title %>)

[Timer]
OnCalendar=<%= @calendar %>
Persistent=true

[Install]
WantedBy=timers.target
