# modules/coal/templates/initscripts/coal.upstart.erb
# Upstart job definition for coal

description "Navigation Timing Whisper logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid coal
setgid coal

respawn
respawn limit unlimited

exec /usr/local/bin/coal --whisper-dir "<%= @whisper_dir %>" --brokers "<%= @kafka_brokers %>" --consumer-group "<%= @kafka_consumer_group %>" --topic "<%= @kafka_topic %>"