# modules/eventlogging/files/mwerrors.conf
; Supervisord configuration for 'mwerrors' Ganglia module.
; Managed by Puppet: puppet:///files/eventlogging/mwerrors.conf
; Forward MediaWiki fatals / exceptions to ZeroMQ

[group:mwerrors]
programs = udp2zmq_8423

[program:udp2zmq_8423]
command = udp2zmq 8423
user = eventlogging
