# modules/eventlogging/files/notebook.conf
[program:notebook]
user = ipython
command = ipython notebook
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = NONE
environment = IPYTHONDIR=/var/lib/ipython
