# modules/limn/templates/init-limn.conf.erb
# Note: This file is managed by Puppet.

# Limn <%= name %> Server Instance

description "Limn <%= name %> Server Instance"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

respawn
respawn limit 2 5

umask 007

kill timeout 300

env NODE_ENV=<%= environment %>
env LIMN_PORT=<%= port %>
env LIMN_VARDIR=<%= var_directory %>

setuid <%= scope.lookupvar('limn::user')  %>
setgid <%= scope.lookupvar('limn::group') %>

console output
chdir <%= base_directory %>
exec ./node_modules/.bin/coco ./server/server.co >> <%= log_file %> 2>&1
