# modules/sentry/templates/sentry.service.erb
[Unit]
Description=Sentry error logging and aggregation service
Documentation=http://sentry.readthedocs.org/en/latest/

[Service]
Type=forking
ExecStart=<%= @sentry_cli %> start
User=<%= @user %>
Restart=always
RestartSec=30

