# modules/service/templates/initscripts/node.upstart.erb
# NOTE: This file is managed by Puppet
# Upstart job for <%= @title %>

description "<%= @title %> service"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

# up ulimit -n a bit
limit nofile <%= @no_file %> <%= @no_file %>

setuid "<%= @title %>"
setgid "<%= @title %>"

env NODE_PATH="/srv/deployment/<%= @repo %>/node_modules"
env <%= @title.gsub(/[^a-zA-Z0-9_]/, '_').upcase %>_PORT="<%= @port %>"

respawn

# wait 60 seconds for a graceful restart before killing the master
kill timeout 60

chdir /srv/deployment/<%= @repo %>
<% if @firejail -%>
exec /usr/bin/firejail --blacklist=root --blacklist=/home/* --tmpfs=/tmp --caps --seccomp /usr/bin/nodejs <%= @starter_script %> -c /etc/<%= @title %>/config.yaml
<% else -%>
exec /usr/bin/nodejs <%= @starter_script %> -c /etc/<%= @title %>/config.yaml
<% end %>
