# modules/mwprof/files/upstart/init.conf
# Upstart job configuration for mwprof
# This file is managed by Puppet
description "MediaWiki profiler"

start on runlevel [2345] or mwprof.start
stop on runlevel [!2345] or mwprof.stop

task

pre-start script
    mkdir -p /run/mwprof
    chown -R mwprof:mwprof /run/mwprof
end script

script
    start mwprof/collector || status mwprof/collector
end script

# vim: set ft=upstart:
