# files/labs/migrate/rc.local
#!/bin/sh -e
/root/migrate_firstboot.sh
exit 0
