# files/txstatsd/init/init.conf
# txstatsd/init
#
description "Initialize txstatsd"
author "Ori Livneh <ori@wikimedia.org>"

start on runlevel [2345] or txstatsd.start

task

script
  for config in /etc/txstatsd/*.cfg; do
    [ -e "$config" ] || break
    name="$(basename $config .cfg)"
    start txstatsd/instance NAME="$name" CONFIG="$config" ||
      status txstatsd/instance NAME="$name"
  done
end script
