# files/txstatsd/init/init.conf
# txstatsd/init
#
# This is an Upstart job configuration file, describing a task job
# (<upstart.ubuntu.com/cookbook/#task-job>). Its role is to start
# service jobs representing the different components of an txstatsd
# server An txstatsd server receives raw log data via UDP and
# transforms it into a sanitized & validated ZeroMQ event stream.
#
description "Initialize txstatsd"
author "Ori Livneh <ori@wikimedia.org>"

start on runlevel [2345] or txstatsd.start

task

script
  for config in /etc/txstatsd/*.cfg; do
    [ -e "$config" ] || break
    [ "$config" = "/etc/txstatsd/txstatsd.conf" ] && continue
    name="$(basename $config .cfg)"
    start txstatsd-multi/instance NAME="$name" CONFIG="$config" ||
      status txstatsd-multi/instance NAME="$name"
  done
end script
