# files/labsnfs/nfs-noidmap.conf
# nfs-noidmap - disable idmap for cross-project NFS
#
# This task disables the kernel's automatic idmap for NFSv4
# to allow conflicting group names accross projects

description "disable idmap for cross-project NFS"

start on starting autofs

task

exec echo 1 >/sys/module/nfs/parameters/nfs4_disable_idmapping
