# modules/mw-rc-irc/files/upstart/ircecho.conf
description "start and stop the irc echo server"
version "1"
author "Chase Pettet"

setuid irc
setgid irc

respawn limit 15 5

start on startup
stop on shutdown

console log

exec /usr/local/bin/udpmxircecho.py rc-pmtpa localhost
