# templates/base/puppet.cron.erb
## This file is managed by Puppet!
## templates/base/puppet.cron.erb
<% $crontime = scope.lookupvar('base::puppet::crontime') -%>
<%= $crontime %>,<%= $crontime.to_i + 30 %>	*	*	*	*	root	timeout <% if scope.function_versioncmp(lsbdistrelease, "12.04") >= 0 %> -k 300<% end %> 1800 puppet agent --onetime --verbose --no-daemonize --no-splay --show_diff >> /var/log/puppet.log 2>&1
@reboot	root	timeout <% if scope.function_versioncmp(lsbdistrelease, "12.04") >= 0 %> -k 300<% end %> 1800 puppet agent --onetime --verbose --no-daemonize --no-splay --show_diff >> /var/log/puppet.log 2>&1
