# modules/librenms/templates/initscripts/librenms-syslog.upstart.erb
# This file is managed by Puppet, librenms module

description "LibreNMS syslog daemon"

start on filesystem
stop on runlevel [06]

expect fork
respawn

script
    exec rsyslogd -c5 -f /etc/librenms-rsyslog.conf -i /var/run/librenms-rsyslogd.pid
end script
