# modules/xdummy/templates/initscripts/xdummy.upstart.erb
# Upstart service job for Xdummy
#
description "Virtual X11 display"

start on runlevel [2345]
stop on runlevel [!2345]

setuid xdummy
setgid xdummy

respawn
respawn limit 10 5

exec Xorg -config /etc/xpra/xorg.conf :<%= @display %>
