# modules/mediawiki/files/hhvm-warmup.conf
# hhvm-warmup - Warm up HHVM's JIT
#
# This task job runs when HHVM starts up. Its purpose is to warm up
# the trace cache by repeatedly requesting URLs. It reads URLs from
# /etc/hhvm/warmup.urls.

description "HHVM warm-up task"
start on started hhvm

task

script
  # Bail if there are no URLs.
  grep -qs '^[^#]' /etc/hhvm/warmup.urls || exit 0

  # Get the JIT warmup request threshold, so we can exceed it by one.
  reqs=$(/usr/bin/hhvm --php -c /etc/hhvm/fcgi.ini -r 'echo ini_get("hhvm.jit_warmup_requests")?:11;')
  reqs=$((reqs+1))

  while read -r uri; do
    # Skip comments and empty lines.
    case "$uri" in "") continue ;; \#*) continue ;; esac
    for _ in $(seq $reqs); do
      /usr/local/bin/furl --headers "$uri" >/dev/null 2>&1
    done
  done < /etc/hhvm/warmup.urls
end script

# vim: set ft=upstart:
