# modules/mediawiki/files/jobrunner.conf
# jobrunner -- Continuously process a MediaWiki job queue
description "MediaWiki job queue loop"

start on (runlevel [2345] and net-device-up IFACE!=lo)
stop on runlevel [!2345]

pre-start script
  mkdir -p -m0755 /var/run/jobrunner
  [ ! -r /etc/default/jobrunner ] && { stop; exit 0; }
  . /etc/default/jobrunner
  [ ! -r "${JOBRUNNER_CONFIG:=/etc/jobrunner.ini}" ] && { stop; exit 1; }
  php -r "exit( !@is_array( parse_ini_file( '$JOBRUNNER_CONFIG' ) ) );" >/dev/null 2>&1
  [ $? -eq 0 ] || { stop; exit 1; }
end script

script
  . /etc/default/jobrunner
  exec start-stop-daemon --quiet --start --pidfile "${JOBRUNNER_PID:-/var/run/jobrunner/pid}" \
    --chuid "${JOBRUNNER_USER:-nobody}" --make-pidfile --exec /usr/local/bin/php -- \
    /srv/deployment/jobrunner/jobrunner/redisJobRunnerService \
    --config-file="${JOBRUNNER_CONFIG:-/etc/jobrunner.ini}" $DAEMON_OPTS
end script

respawn

# vim: set ft=upstart:
