# modules/ganglia-new/files/upstart/ganglia-monitor.conf
# files/ganglia-monitor.conf

description "Ganglia Monitor daemon"
author "Mark Bergsma <mark@wikimedia.org>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec /usr/sbin/gmond