# templates/udp2log/wfdebug-ganglia.conf.erb
# vim: set ft=upstart:

# Upstart job configuration for wfdebug-ganglia.py
# This file is managed by Puppet
description "Report MediaWiki log counts to Ganglia"
start on (local-filesystems and net-device-up IFACE!=lo)
setuid nobody
exec python /usr/local/bin/wfdebug-ganglia.py <%= @ganglia_reporter_port %>
respawn
