# modules/eventlogging/files/init/init.conf
# eventlogging/init
#
# This is an Upstart job configuration file, describing a task job
# (<upstart.ubuntu.com/cookbook/#task-job>). Its role is to start
# service jobs representing the different components of an EventLogging
# server An EventLogging server receives raw log data via UDP and
# transforms it into a sanitized & validated ZeroMQ event stream.
#
description "Initialize EventLogging"
author "Ori Livneh <ori@wikimedia.org>"

start on runlevel [2345] or eventlogging.start

task

script
    roles="forwarder processor multiplexer consumer reporter"
    set -- $roles
    for role in "$@"; do
        for config in /etc/eventlogging.d/${role}s/*; do
            [ -e "$config" ] || break
            name="$(basename $config)"
            start eventlogging/$role NAME="$name" CONFIG="$config" ||
                status eventlogging/$role NAME="$name"
        done
    done
end script
