# modules/generic/files/upstart/enable-rps.conf
# enable-rps

description "Enable RPS on eth0 receive queues"
author "Mark Bergsma <mark@wikimedia.org>"

start on filesystem
task
script
	for queue in /sys/class/net/eth0/queues/rx-*; do echo ff > $queue/rps_cpus; done
end script
