# modules/generic/files/upstart/fcgi-thumb-handler.conf
# This file will be managed by Puppet!
#
# fcgi-thumb-handler

description	"Thumbnail generation FastCGI handler"

start on filesystem
stop on runlevel S

respawn
respawn limit 10 5

umask 022

exec /usr/bin/spawn-fcgi -n -u 127.0.0.1 -p 9000 -u www-data -C 60 -- /usr/bin/php5-cgi
