# modules/generic/files/upstart/manage-nfs-volumes.conf
description	"Manage shared NFS volumes for Labs projects"

start on runlevel [2345]

respawn
respawn limit 10 5

exec su -c "/usr/bin/python /usr/local/sbin/manage-nfs-volumes-daemon --logfile=/var/lib/nfsmanager/manage-nfs-volumes.log" nfsmanager
