# modules/generic/files/upstart/manage-volumes.conf
description	"Manage shared gluster volumes for Labs projects"

start on runlevel [2345]

respawn
respawn limit 10 5

exec su -c "/usr/bin/python /usr/local/sbin/manage-volumes-daemon --logfile=/var/lib/glustermanager/manage-volumes.log" glustermanager
