# modules/generic/files/upstart/nfs-noidmap.conf
# nfs-noidmap - disable idmap for cross-project NFS
#
# This task disables the kernel's automatic idmap for NFSv4
# to allow conflicting group names accross projects
#
# Made into an upstart task because that setting is annoyingly
# not made available through sysctl.

description "disable idmap for cross-project NFS"

start on starting autofs

task

exec echo 1 >/sys/module/nfs/parameters/nfs4_disable_idmapping
