# modules/generic/files/upstart/twemproxy.conf
# vim: set ft=upstart:

# Upstart job configuration for twemproxy
# This file is managed by Puppet

description "twemproxy"

start on (net-device-up
        and runlevel [2345])
stop on runlevel [!2345]
setuid nobody
respawn
respawn limit 5 5

pre-start script
        # Cannot use mw-deployment-vars.sh as not compatible with dash
        # [ -f /usr/local/lib/mw-deployment-vars.sh ] || { stop; exit 0; }
        # . /usr/local/lib/mw-deployment-vars.sh

        env MW_COMMON = "/usr/local/apache/common-local"
        [ -f "$MW_COMMON/multiversion/MWRealm.sh" ] || { stop; exit 0; }
        . $MW_COMMON/multiversion/MWRealm.sh
        env config = `getRealmSpecificFilename "$MW_COMMON_SOURCE/wmf-config/twemproxy.yaml"`

        [ -f "$config" ] || { stop; exit 0; }
end script

exec /usr/local/bin/nutcracker -m 65536 -a 127.0.0.1 -c "$config"
