# modules/generic/files/upstart/twemproxy.conf
description "twemproxy"

start on (net-device-up
        and runlevel [2345])
stop on runlevel [!2345]
setuid nobody
respawn
respawn limit 5 5

pre-start script
        [ -f /etc/wikimedia-site ] || { stop; exit 0; }
        [ -f "/usr/local/apache/common/wmf-config/twemproxy-$(cat /etc/wikimedia-site)".yaml ] || { stop; exit 0; }
end script

exec /usr/local/bin/nutcracker -m 65536 -a 127.0.0.1 -c "/usr/local/apache/common/wmf-config/twemproxy-$(cat /etc/wikimedia-site).yaml"

