# modules/generic/files/upstart/twemproxy.conf
# vim: set ft=upstart:

# Upstart job configuration for twemproxy
# This file is managed by Puppet

description "twemproxy"

start on (net-device-up
        and runlevel [2345])
stop on runlevel [!2345]
setuid nobody
respawn
respawn limit 5 5

pre-start script
        [ -f /usr/local/lib/mw-deployment-vars.sh ] || { stop; exit 0; }
        . /usr/local/lib/mw-deployment-vars.sh

        [ -f "$MW_COMMON_SOURCE/multiversion/MWRealm.sh" ] || { stop; exit 0; }
        . $MW_COMMON_SOURCE/multiversion/MWRealm.sh
        env config = `getRealmSpecificFilename "$MW_COMMON_SOURCE/wmf-config/twemproxy.yaml"`

        [ -f "/usr/local/apache/common/wmf-config/$config" ] || { stop; exit 0; }
end script

exec /usr/local/bin/nutcracker -m 65536 -a 127.0.0.1 -c "/usr/local/apache/common/wmf-config/$config"

