# modules/brrd/templates/brrd.conf.erb
# Upstart job definition for brrd, an RRD metric logger.
description "RRD metric logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid brrd
setgid brrd

respawn
respawn limit unlimited

chdir /srv/deployment/brrd/brrd
exec /usr/bin/python -m brrd            \
    --log-file "/var/log/brrd/brrd.log" \
    "<%= @endpoint %>"                  \
    "/var/lib/brrd/navtiming.rrd"
