# modules/ganglia/files/upstart/ganglia-monitor-aggregator.conf
# files/ganglia-monitor-aggregator.conf

description "Ganglia Monitor aggregator"
author "Mark Bergsma <mark@wikimedia.org>"

start on runlevel [2345]
stop on runlevel [!2345]

task

script
	for gmonid in $(ls /etc/ganglia/aggregators/*.conf)
	do
		start ganglia-monitor-aggregator-instance ID=${gmonid%.conf}
	done
end script
