# modules/xvfb/templates/initscripts/xvfb.upstart.erb
# Upstart service job for xvfb
#
description "X virtual framebuffer"

start on runlevel [2345]
stop on runlevel [!2345]

setuid xvfb
setgid xvfb

respawn
respawn limit 10 5

exec Xvfb :<%= @display %> -screen 0 <%= @resolution %> -ac -noreset
