# modules/uwsgi/templates/initscripts/uwsgi-startup.upstart.erb
description "create required directories for sockets/stats"

start on mounted MOUNTPOINT=/run TYPE=tmpfs

task

exec install -d -o root -g root /run/uwsgi
