# modules/coal/templates/coal.conf.erb
# Upstart job definition for coal

description "Navigation Timing Whisper logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid coal
setgid coal

respawn
respawn limit unlimited

exec /usr/local/bin/coal --whisper-dir "/var/lib/coal" "<%= @endpoint %>"
