# modules/graphite/templates/graphite-web.conf.erb
# graphite-web
#
# Upstart job configuration for Graphite webapp.
# Part of the Graphite suite.
#
description "Graphite web app"

start on started memcached

setuid www-data
setgid www-data

chdir /usr/share/pyshared/graphite

exec /usr/bin/uwsgi_python \
    --socket /var/run/graphite-web/graphite-web.sock \
    --wsgi-file /usr/share/graphite-web/graphite.wsgi \
    --die-on-term \
    --master \
    --processes <%= @uwsgi_processes %>

respawn
