# modules/base/templates/puppet.cron.erb
#####################################################################
##### THIS FILE IS MANAGED BY PUPPET
#####  as template('base/puppet.cron.erb')
######################################################################
<%-
interval = scope.lookupvar('base::puppet::interval').to_i
crontime = scope.lookupvar('base::puppet::crontime').to_i
numtimes = 60 / interval
tmp = Array.new(numtimes) { |t| t = t * interval + crontime }
tmp = tmp.map { |x| if x < 60 then x  else x - 60 end }
tmp = tmp.sort()
times = tmp.join(',')
-%>
<%= times %> * * * * root /usr/local/sbin/puppet-run > /dev/null 2>&1
@reboot              root /usr/local/sbin/puppet-run > /dev/null 2>&1
