# modules/graphite/templates/initscripts/frontend-relay.upstart.erb
description "carbon-c-relay instance, used to route metrics"

start on stopped carbon/init
stop on runlevel [!2345] or carbon.stop

respawn

limit nofile 32768 32768

setuid carbon-c-relay

# it is possible there's another carbon-c-relay running on the same host, so
# change the hostname
exec /usr/bin/carbon-c-relay -H <%= "#{hostname}_frontend" -%> -p 2003 -f /etc/carbon/frontend-relay.conf
