# modules/tc/templates/initscripts/tc.upstart.erb
description "Apply traffic control rules from /etc/tc/*.d"

start on started networking

task

exec /usr/local/sbin/apply-tc
