# modules/tc/templates/initscripts/tc.upstart.erb
description "Apply traffic control rules from /etc/tc"

start on started networking

task

exec /bin/run-parts /etc/tc
