# modules/eventlogging/templates/mwerrors.conf.erb
; Supervisord configuration for 'mwerrors' Ganglia module.
; File managed by Puppet
; Forward MediaWiki fatals / exceptions to ZeroMQ

[group:mwerrors]
programs = udp2zmq_<%= @port %>

[program:udp2zmq_<%= @port %>]
command = udp2zmq <%= @port %>
user = eventlogging
