# modules/eventlogging/templates/notebook.conf.erb
; supervisor configuration for ipython-notebook
; This file is managed by Puppet.

[program:notebook]
command = ipython notebook --ipython-dir="<%= ipython_dir %>" --profile="<%= ipython_profile %>"
user = <%= ipython_user %>
directory = <%= ipython_dir %>
redirect_stderr = true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = NONE
environment = IPYTHONDIR=<%= ipython_dir %>,MPLCONFIGDIR=<%= graph_dir %>
