# modules/bird/templates/anycast-healthchecker.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=Anycast healthchecker
After=network.target
Requires=network.target
Documentation=https://github.com/unixsurfer/anycast_healthchecker/blob/master/README.rst
<%- if @bind_service -%>
After=<%= @bind_service %>
BindsTo=<%= @bind_service %>
<%- end -%>

[Service]
Type=simple
TimeoutStartSec=0
User=bird
Group=bird
ExecStart=/usr/bin/anycast-healthchecker
Restart=on-failure

[Install]
WantedBy=multi-user.target
# Note no wantedby on the bind_service here, because we should be able to keep
# the healthechecker and bird explicitly-stopped while cycling the liveness of
# the underlying service manually.
