# modules/bird/templates/bird.service.erb
[Unit]
Description=BIRD Internet Routing Daemon (IPv4)

<%- if @bind_service -%>
After=<%= @bind_service %>
BindsTo=<%= @bind_service %>
<%- end -%>

[Service]
EnvironmentFile=/etc/bird/envvars
ExecStartPre=/usr/lib/bird/prepare-environment
ExecStartPre=/usr/sbin/bird -p
ExecReload=/usr/sbin/birdc configure
ExecStart=/usr/sbin/bird -f -u $BIRD_RUN_USER -g $BIRD_RUN_GROUP $BIRD_ARGS
Restart=on-abort

[Install]
<%- if @bind_service -%>
## https://github.com/systemd/systemd/issues/720
WantedBy=<%= @bind_service %>
<%- end -%>
