# modules/dynamicproxy/files/upstart.conf
description	"API for the dynamic proxy, aka Invisible Unicorn"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

script
    /usr/bin/uwsgi --logto /var/log/uwsgi/app/invisible.log --buffer-size=32768 --plugin python --die-on-term -s /tmp/uwsgi.sock --module invisible_unicorn.api --callable app --processes 4 --threads 2 --master true --uid www-data
end script
