# modules/profile/templates/mjolnir/kafka-bulk-daemon.service.erb
[Unit]
Description=MjoLniR kafka bulk daemon
After=network.target

[Service]
SyslogIdentifier=mjolnir-kafka-bulk-daemon
User=nobody
Group=nogroup
Environment=REQUESTS_CA_BUNDLE=/etc/ssl/certs/ca-certificates.crt
ExecStart=/srv/deployment/search/mjolnir/deploy/venv/bin/mjolnir-utilities.py kafka_bulk_daemon \
        --brokers <%= @kafka_config['brokers']['string'] %> \
        --es-clusters <%= @es_cluster_endpoint %> \
        --group-id <%= @group_id %> \
        <%- if @topics.any? -%>
        --topic <%= @topics.join(' ') %> \
        <%- end -%>
        <%- if @priority_topics.any? -%>
        --priority-topic <%= @priority_topics.join(' ') %> \
        <%- end -%>
        --prometheus-port <%= @prometheus_port %>

StandardInput=null
RestartSec=2
Restart=always

[Install]
WantedBy=multi-user.target
