# modules/profile/templates/mjolnir/kafka-daemon.service.erb
[Unit]
Description=MjoLniR kafka daemon
After=network.target

[Service]
User=nobody
Group=nogroup
ExecStart=/srv/deployment/search/mjolnir/venv/bin/mjolnir-utilities.py kafka_daemon --brokers <%= @kafka_config['brokers']['string'] %>
StandardInput=null
StandardOutput=journal
StandardError=journal
RestartSec=2

[Install]
WantedBy=multi-user.target
