# modules/profile/templates/mjolnir/kafka-bulk-daemon.service.erb
[Unit]
Description=MjoLniR kafka bulk daemon
After=network.target

[Service]
SyslogIdentifier=mjolnir-kafka-bulk-daemon
User=nobody
Group=nogroup
ExecStart=/srv/deployment/search/mjolnir/venv/bin/mjolnir-utilities.py kafka_bulk_daemon \
        --brokers <%= @kafka_config['brokers']['string'] %> \
        --es-clusters localhost:9200 \
        --group-id <%= @group_id %> \
        --topic <%= @topics.join(' ') %> \
        --prometheus-port <%= @prometheus_port %>

StandardInput=null
RestartSec=2
Restart=always

[Install]
WantedBy=multi-user.target
