# modules/profile/templates/mjolnir/kafka-msearch-daemon.service.erb
[Unit]
Description=MjoLniR kafka msearch daemon
After=network.target

[Service]
SyslogIdentifier=mjolnir-kafka-msearch-daemon
User=nobody
Group=nogroup
ExecStart=/srv/deployment/search/mjolnir/venv/bin/mjolnir-utilities.py kafka_msearch_daemon \
    --brokers <%= @kafka_config['brokers']['string'] %> \
    --input-topic <%= @input_topic %> \
    --output-topic <%= @output_topic %> \
    --num-workers <%= @num_workers %> \
    --max-concurrent-searches <%= @max_concurrent_searches %> \
    --prometheus-port <%= @prometheus_port %>

StandardInput=null
RestartSec=2
Restart=always

[Install]
WantedBy=multi-user.target
