# modules/cfssl/templates/cfssl-ocspserve.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=Cloudflare SSL OCSP Responder (<%= @title %>)
After=network.target remote-fs.target nss-lookup.target
Documentation=https://github.com/cloudflare/cfssl/tree/master/doc

[Service]
ExecStart=/usr/bin/cfssl ocspserve \
          -address <%= @listen_addr %> \
          -port <%= @listen_port %> \
          -responses <%= @_responses_file %>
Restart=always
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
