# modules/cfssl/templates/cfssl.service.erb
[Unit]
Description=Cloudflare SSL
After=network.target remote-fs.target nss-lookup.target
Documentation=https://github.com/cloudflare/cfssl/tree/master/doc

[Service]
ExecStart=/usr/bin/cfssl serve \
          -db-config <%= @db_conf_file %> \
          -config <%= @conf_file %> \
          -address <%= @host %> \
          -port <%= @port %> \
          -ca <%= @ca_file %> \
          -ca-key <%= @ca_key_file %>
Restart=always
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
