# modules/cfssl/templates/cfssl.service.erb
[Unit]
Description=Cloudflare SSL (<%= @title %>)
After=network.target remote-fs.target nss-lookup.target
Documentation=https://github.com/cloudflare/cfssl/tree/master/doc

[Service]
ExecStart=/usr/bin/cfssl serve \
          -db-config <%= @_db_conf_file %> \
          -config <%= @conf_file %> \
          -address <%= @host %> \
          -port <%= @port %> \
<%- if @ca_bundle_file -%>
          -ca-bundle <%= @ca_bundle_file %> \
<%- end -%>
<%- if @enable_ocsp -%>
          -responder <%= @ocsp_cert_path %> \
          -responder-key <%= @ocsp_key_path %> \
<%- end -%>
          -ca <%= @_ca_file %> \
          -ca-key <%= @_ca_key_file %>
Restart=always
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
