# modules/keyholder/templates/keyholder-agent.conf.erb
# keyholder-agent - Shared SSH-agent
#
# Runs the ssh-agent(1) instance that holds shared identities.

description "Shared SSH agent"

start on (local-filesystems and net-device-up IFACE!=lo)

setgid keyholder
setuid keyholder

exec /usr/bin/ssh-agent -d -a <%= @agent_socket %>
post-start exec [ -S <%= @agent_socket %> ] || sleep 1
post-stop exec /bin/rm -f <%= @agent_socket %>

# vim: set ft=upstart:
