# modules/keyholder/templates/keyholder-proxy.conf.erb
# keyholder-proxy - Filtering proxy for ssh-agent(1)
#
# The `keyholder-proxy` service runs the filtering ssh-agent proxy
# that acts as an intermediary between users in the trusted group
# and the backend ssh-agent that holds the shared key(s).

description "Shared SSH agent proxy"

start on started keyholder-agent
stop on stopped keyholder-agent

setuid keyholder
setgid <%= @trusted_group %>

umask 007

pre-start exec /bin/rm -f <%= @proxy_socket %>
exec /usr/local/bin/ssh-agent-proxy --bind <%= @proxy_socket %> --connect <%= @agent_socket %>
post-stop exec /bin/rm -f <%= @proxy_socket %>

# vim: set ft=upstart:
