# modules/txstatsd/files/txstatsd.conf
# vim: set ft=upstart:
# Upstart job definition for txStatsD
description "txStatsD metrics aggregation network daemon"

setuid txstatsd
setgid txstatsd

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec twistd -n -l- --pidfile= statsd \
    --config=/etc/txstatsd/txstatsd.cfg
