# modules/gitlab_runner/templates/gitlab-runner.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Service]
<%- unless @gitlab_runner_user == 'root' -%>
User=
User=<%= @gitlab_runner_user %>
ExecStart=
ExecStart=/usr/bin/gitlab-runner "run" "--working-directory" "/home/<%= @gitlab_runner_user %>" "--config" "<%= @config_path%>/config.toml" "--service" "gitlab-runner" "--user" "<%= @gitlab_runner_user %>"
<% end -%>