# modules/gitlab_runner/templates/gitlab-runner.service.erb
[Service]
User=
User=<%= @gitlab_runner_user %>