# modules/chromium/templates/chromium.conf.erb
# chromium - Chromium web browser
#
# Chromium is the open source web browser project from which Google
# Chrome draws its source code.
#
description "Chromium web browser"

start on started xvfb

setuid chromium
setgid chromium

respawn
respawn limit 10 5

exec chromium-browser \
    <% if @incognito %>--incognito<% end %> \
    --remote-debugging-port=<%= @remote_debugging_port %> \
    --display=:<%= @display.to_s.sub(/\D/, '') %> \
    <%= Array(@extra_args).join(" \\\n    ") %>
