# modules/rcstream/files/upstart/init.conf
# rcstream/init
#
# This is an Upstart job configuration file, describing a task job.
# It sources /etc/default/rcstream and starts an RCStream server
# instance for each configured port.
#
description "RCStream master task"

start on runlevel [2345] or rcstream.start
task

script
  [ -r /etc/default/rcstream ] || { exit; }
  . /etc/default/rcstream
  IFS=,
  set -- $ports
  for port in "$@"; do
    server=$iface:$port
    start rcstream/server server=$server redis=$redis ||
      status rcstream/server server=$server
  done
end script
