# modules/mw_rc_irc/files/systemd/ircd.service
[Unit]
Description=IRCd for Mediawiki RecentChanges feed
After=network.target

[Service]
User=irc
Group=irc
WorkingDirectory=/var/lib/gitblit
ExecStart=/usr/bin/ircd -foreground
[Install]
WantedBy=multi-user.target
