# modules/buildkitd/templates/buildkitd.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=buildkitd
After=docker.service
Requires=docker.service

[Service]
User=buildkitd
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker rm -f buildkitd
ExecStart=/usr/bin/docker run --name buildkitd \
    --security-opt seccomp=unconfined --security-opt apparmor=unconfined \
    --network=<%= @network %> \
    <%= @image %> \
    --addr tcp://<%= @address %>:<%= @port %>
ExecStop=/usr/bin/docker stop buildkitd
Restart=on-failure

[Install]
WantedBy=multi-user.target
