# modules/buildkitd/templates/buildkitd.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=buildkitd
After=docker.service
Requires=docker.service

[Service]
User=buildkitd
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker rm -f buildkitd
ExecStart=/usr/bin/docker run --name buildkitd \
    --security-opt seccomp=unconfined --security-opt apparmor=unconfined \
    --volume /etc/buildkitd.toml:/etc/buildkit/buildkitd.toml \
    --network <%= @network %> \
    <%- if @enable_webproxy -%>
    --env http_proxy=<%= @http_proxy %> \
    --env https_proxy=<%= @https_proxy %> \
    --env no_proxy=<%= @no_proxy %> \
    <%- end -%>
    <%= @image %> \
    --addr tcp://<%= @address %>:<%= @port %> \
    --config /etc/buildkit/buildkitd.toml
ExecStop=/usr/bin/docker stop buildkitd
Restart=on-failure

[Install]
WantedBy=multi-user.target
