# modules/changes/files/upstart/init.conf
# changes/init
#
# This is an Upstart job configuration file, describing a task job.
# It sources /etc/default/changes and starts an Changes server
# instance for each configured port.
#
description "Changes master task"

start on runlevel [2345] or changes.start
task

script
  [ -r /etc/default/changes ] || { exit; }
  . /etc/default/changes
  IFS=,
  set -- $ports
  for port in "$@"; do
    server=$iface:$port
    start changes/server server=$server redis=$redis ||
      status changes/server server=$server
  done
end script
