# modules/confd/templates/initscripts/confd.upstart.erb
# confd config management daemon
description "Confd config management daemon"

start on (local-filesystem and net-device-up IFACE!=lo)

script
    CONFD_BACKEND="<%= @backend %>"
    CONFD_OPTS="-node <%= @node %> <% if @interval %>-interval <%= @interval %><% else %>-watch<% end %>"
    /usr/bin/confd -backend $CONFD_BACKEND $CONFD_OPTS
end script

# keep this up, but do not exceed 10 attempts at restart in 30 seconds
respawn
respawn limit 10 30
