# modules/openstack/templates/xena/neutron/neutron-api
#! /bin/sh
# SPDX-License-Identifier: Apache-2.0
### BEGIN INIT INFO
# Provides:          neutron-api
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      mysql postgresql rabbitmq-server keystone
# Should-Stop:       mysql postgresql rabbitmq-server keystone
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: neutron-api
# Description:       Provides the Neutron networking API
### END INIT INFO

DESC="OpenStack Neutron API (neutron-api)"
PROJECT_NAME=neutron
NAME=${PROJECT_NAME}-api
UWSGI_PORT=<%= @bind_port %>
UWSGI_INI_PATH=/etc/neutron/neutron-api-uwsgi.ini
UWSGI_INI_APP=/usr/bin/neutron-api

### Maintain the plugin selection so that we can load the corresponding .ini file ###
PYARGV="--config-file=/etc/neutron/neutron.conf"
[ -r /usr/share/neutron-common/plugin_guess_func ] || exit 0
. /usr/share/neutron-common/plugin_guess_func

if ! [ -r /etc/neutron/neutron.conf ] ; then
	echo "Cloud not read /etc/neutron/neutron.conf: exiting"
	exit 0
fi

CURRENT_PLUGIN=`grep "^[ \t]*core_plugin[ \t]*=[ \t]*[._a-zA-Z0-9]*\$" /etc/neutron/neutron.conf | sed -e 's/^[ \t]*core_plugin[ \t]*=[ \t]*//'`
if [ -z "${CURRENT_PLUGIN}" ] ; then
	echo "No core_plugin= value found: please set it and try again"
	exit 0
fi
neutron_core_plugin_to_plugin_name ${CURRENT_PLUGIN}
neutron_plugin_ini_path ${NEUTRON_PLUGIN_NAME}
if [ -z "${NEUTRON_PLUGIN_CONFIG}" ] ; then
	echo "Plugin not recognized: please edit /etc/init.d/neutron-api to select the correct .ini file to load for your plugin"
else
	PYARGV="${PYARGV} --config-file=${NEUTRON_PLUGIN_CONFIG}"
	DESC="${DESC} with ${NEUTRON_PLUGIN_NAME} plugin"
fi

# Load server side config files
if [ -e /etc/neutron/server.conf.d ] ; then
	PYARGV="--config-dir=/etc/neutron/server.conf.d ${PYARGV}"
fi

### Neutron folders creation ###
mkdir -p /var/lib/neutron/tmp
chown neutron:neutron /var/lib/neutron/tmp
export TMPDIR=/var/lib/neutron/tmp
#!/bin/sh
# The content after this line comes from openstack-pkg-tools
# and has been automatically added to a .init.in script, which
# contains only the descriptive part for the daemon. Everything
# else is standardized as a single unique script.

# Author: Thomas Goirand <zigo@debian.org>
# Author: Ondřej Nový <novy@ondrej.org>

# PATH should only include /usr/* if it runs after the mountnfs.sh script
PATH=/sbin:/usr/sbin:/bin:/usr/bin

if [ -n "${UWSGI_PORT}" ] && [ -n "${UWSGI_INI_PATH}" ] && [ -n "${UWSGI_INI_APP}" ] ; then
	if ! [ -f "${UWSGI_INI_APP}" ] ; then
		exit 0
	fi
	if [ -d /etc/${PROJECT_NAME}/ssl/private ] ; then
		KEY_FILE=$(find /etc/${PROJECT_NAME}/ssl/private -type f -iname '*.pem' 2>/dev/null | head -n 1)
	fi

	if [ -e /usr/local/share/ca-certificates/puppet_openstack.crt ] ; then
		# This is needed for puppet...
		CERT_FILE=/usr/local/share/ca-certificates/puppet_openstack.crt
	else
		if [ -d /etc/${PROJECT_NAME}/ssl/public ] ; then
			CERT_FILE=$(find /etc/${PROJECT_NAME}/ssl/public -type f -iname '*.crt' 2>/dev/null | head -n 1)
		fi
	fi

	# Sid doesn't have /usr/bin/uwsgi_python3, so we need
	# to search for a more specific daemon name. For stretch
	# /usr/bin/uwsgi_python3 is fine.
	for i in 3 35 36 37 38 39 ; do
		if [ -x /usr/bin/uwsgi_python${i} ] ; then
			DAEMON=/usr/bin/uwsgi_python${i}
		fi
	done

	if [ -n "${KEY_FILE}" ] && [ -n "${CERT_FILE}" ] ; then
		DAEMON_ARGS="--https-socket [::]:${UWSGI_PORT},${CERT_FILE},${KEY_FILE}"
	else
		DAEMON_ARGS="--http-socket [::]:${UWSGI_PORT}"
	fi

	DAEMON_ARGS="${DAEMON_ARGS} --ini ${UWSGI_INI_PATH}"
	NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
	NO_OPENSTACK_LOGFILE_DAEMON_ARG=yes
fi

if [ -z "${DAEMON}" ] ; then
	DAEMON=/usr/bin/${NAME}
fi
PIDFILE=/var/run/${PROJECT_NAME}/${NAME}.pid
if [ -z "${SCRIPTNAME}" ] ; then
	SCRIPTNAME=/etc/init.d/${NAME}
fi
if [ -z "${SYSTEM_USER}" ] ; then
	SYSTEM_USER=${PROJECT_NAME}
fi
if [ -z "${SYSTEM_GROUP}" ] ; then
	SYSTEM_GROUP=${PROJECT_NAME}
fi
if [ "${SYSTEM_USER}" != "root" ] ; then
	STARTDAEMON_CHUID="--chuid ${SYSTEM_USER}:${SYSTEM_GROUP}"
fi
if [ -z "${CONFIG_FILE}" ] ; then
	CONFIG_FILE=/etc/${PROJECT_NAME}/${PROJECT_NAME}.conf
fi
LOGFILE=/var/log/${PROJECT_NAME}/${NAME}.log
if [ -z "${NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG}" ] ; then
	DAEMON_ARGS="--config-file=${CONFIG_FILE} ${DAEMON_ARGS}"
fi

# Exit if the package is not installed
[ -x $DAEMON ] || exit 0

# If ran as root, create /var/lock/X, /var/run/X and /var/cache/X as needed
if [ `whoami` = "root" ] ; then
	for i in lock run cache ; do
		mkdir -p /var/$i/${PROJECT_NAME}
		chown ${SYSTEM_USER}:${SYSTEM_GROUP} /var/$i/${PROJECT_NAME}
	done
fi

# This defines support functions which we use later on
. /lib/lsb/init-functions

RET=0

# Force logging to /var/log/neutron/neutron-api.log
#
#  The upstream script checks defaults but we
#  always want this.
DAEMON_ARGS="$DAEMON_ARGS --logto2=$LOGFILE"

do_start() {
	start-stop-daemon \
		--start \
		--quiet \
		--background ${STARTDAEMON_CHUID} \
		--make-pidfile --pidfile ${PIDFILE} \
		--chdir /var/lib/${PROJECT_NAME} \
		--startas $DAEMON \
		--test > /dev/null \
		|| return 1
	if [ -n "${PYARGV}" ] ; then
		start-stop-daemon \
			--start \
			--quiet \
			--background ${STARTDAEMON_CHUID} \
			--make-pidfile --pidfile ${PIDFILE} \
			--chdir /var/lib/${PROJECT_NAME} \
			--startas $DAEMON \
			-- $DAEMON_ARGS --pyargv "${PYARGV}" \
			|| return 2
	else
		start-stop-daemon \
			--start \
			--quiet \
			--background ${STARTDAEMON_CHUID} \
			--make-pidfile --pidfile ${PIDFILE} \
			--chdir /var/lib/${PROJECT_NAME} \
			--startas $DAEMON \
			-- $DAEMON_ARGS \
			|| return 2
	fi
}

do_stop() {
	start-stop-daemon \
		--stop \
		--quiet \
		--retry=TERM/30/KILL/5 \
		--pidfile $PIDFILE
	RETVAL=$?
	rm -f $PIDFILE
	return "$RETVAL"
}

do_systemd_start() {
	if [ -n "${PYARGV}" ] ; then
		exec $DAEMON $DAEMON_ARGS --pyargv "${PYARGV}"
	else
		exec $DAEMON $DAEMON_ARGS
	fi
}

case "$1" in
start)
	log_daemon_msg "Starting $DESC" "$NAME"
	do_start
	case $? in
		0|1) log_end_msg 0 ; RET=$? ;;
		2)   log_end_msg 1 ; RET=$? ;;
	esac
;;
stop)
	log_daemon_msg "Stopping $DESC" "$NAME"
	do_stop
	case $? in
		0|1) log_end_msg 0 ; RET=$? ;;
		2)   log_end_msg 1 ; RET=$? ;;
	esac
;;
status)
	status_of_proc "$DAEMON" "$NAME"
	RET=$?
;;
systemd-start)
	do_systemd_start
;;  
show-args)
	if [ -n "${PYARGV}" ] ; then
		echo $DAEMON $DAEMON_ARGS --pyargv \"${PYARGV}\"
	else
		echo $DAEMON $DAEMON_ARGS
	fi
;;
restart|force-reload)
	log_daemon_msg "Restarting $DESC" "$NAME"
	do_stop
	case $? in
		0|1)
			do_start
			case $? in
				0) log_end_msg 0 ; RET=$? ;;
				1) log_end_msg 1 ; RET=$? ;; # Old process is still running
				*) log_end_msg 1 ; RET=$? ;; # Failed to start
			esac
		;;
		*) log_end_msg 1 ; RET=$? ;; # Failed to stop
	esac
;;
*)
	echo "Usage: $SCRIPTNAME {start|stop|status|restart|force-reload|systemd-start}" >&2
	RET=3
;;
esac

exit $RET
