# modules/ganglia_new/templates/gmetad.upstart
# Upstart job for Ganglia gmetad

description "Ganglia Metadata Daemon"
author "Bernard Li <bernard@vanhpc.org> Alexandros Kosiaris <akosiaris@wikimedia.org>"

expect daemon
respawn

start on runlevel [2345]
stop on runlevel [!2345]

env PIDFILE="/var/run/gmetad.pid"
env RRDCACHED_ADDRESS="<%= @rrdcached_socket %>"

exec /usr/sbin/gmetad --pid-file=$PIDFILE
