# modules/graphite/files/carbon-upstart/init.conf
# Upstart job configuration for Carbon, a part of Graphite
# This file is managed by Puppet
description "Graphite metrics storage"

# Start after filesystems are mounted and networking is up
start on (filesystem and net-device-up IFACE!=lo) or carbon.start
stop on runlevel [!2345] or carbon.stop

task

pre-start script
    mkdir -p /var/run/carbon
    chown _graphite:_graphite /var/run/carbon
end script

script
    # Start an instance for each '[cache:<instance name>]' match.
    instances=$(grep -Po '(?<=^\[cache:).+(?=\]$)' /etc/carbon/carbon.conf)
    set -- $instances
    for name in "$@"; do
        start carbon/cache NAME="$name" || status carbon/cache NAME="$name"
    done
end script

# vim: set ft=upstart:
