# modules/graphite/files/carbon-upstart/relay.conf
# graphite/carbon-relay
#
# carbon-relay.py listens for metrics and shards or replicates across multiple
# carbon-cache.py instances. Part of the Graphite suite.
#
description "Graphite metrics listener"

# Run after carbon/init has spun up all carbon/cache instances.
start on stopped carbon/init
stop on runlevel [!2345] or carbon.stop

setuid _graphite
setgid _graphite

pre-start script
    # Exit normally if no relay daemon is configured for this node.
    grep -Fqx '[relay]' /etc/carbon/carbon.conf || { stop; exit 0; }
end script

# "--debug" makes carbon-* stay in the foreground and log to stdout, which
# allows Upstart to take care of daemonizing and logging.
exec /usr/bin/carbon-relay --config="/etc/carbon/carbon.conf" --debug start

respawn
