# modules/webperf/templates/coal.conf.erb
# Upstart job definition for coal

description "Navigation Timing Whisper logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid webperf
setgid webperf

respawn
respawn limit unlimited

exec /srv/webperf/coal --whisper-dir "/var/lib/coal" "<%= @endpoint %>"
