# modules/webperf/templates/navtiming.conf.erb
# navtiming
#
# This is an Upstart job configuration file for a Graphite metric module for
# NavigationTiming events. For more information, see:
# https://meta.wikimedia.org/wiki/Schema:NavigationTiming
# http://www.mediawiki.org/wiki/Extension:NavigationTiming
#
# This file is managed by Puppet.
#
description "NavigationTiming Graphite module"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid webperf
setgid webperf

respawn
respawn limit 15 5

chdir /srv/webperf
exec /usr/bin/python /srv/webperf/navtiming.py \
    <%= @endpoint %> \
    --statsd-host <%= @statsd_host %> \
    --statsd-port <%= @statsd_port %>
