# modules/webperf/templates/rrd-navtiming.conf.erb
# rrd-navtiming

description "NavigationTiming RRD logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid webperf
setgid webperf

respawn
respawn limit unlimited

exec /srv/webperf/rrd-navtiming         \
    --endpoint "<%= @endpoint %>"       \
    --rrd-path "/var/lib/rrd-navtiming" \
    --log-path "/var/log/rrd-navtiming"
