# modules/webperf/templates/statsv.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=statsv
After=network-online.target

[Service]
WorkingDirectory=/srv/deployment/statsv/statsv
ExecStart=/usr/bin/python3 /srv/deployment/statsv/statsv/statsv.py --brokers <%= @kafka_brokers %> --security-protocol <%= @kafka_security_protocol %><%= " --ssl-cafile #{@kafka_ssl_cafile}" if @kafka_ssl_cafile %> --statsd <%= @statsd_host %>:<%= @statsd_port %> --topics <%= @topics %> <%= "--api-version #{@kafka_api_version}" if @kafka_api_version %>
User=nobody
Restart=always
WatchdogSec=120

[Install]
WantedBy=multi-user.target
