# modules/editstream/files/upstart/init.conf
# editstream/init
#
# This is an Upstart job configuration file, describing a task job.
# It sources /etc/default/editstream and starts an EditStream server
# instance for each configured port.
#
description "EditStream master task"

start on runlevel [2345] or editstream.start
task

script
  [ -r /etc/default/editstream ] || { exit; }
  . /etc/default/editstream
  IFS=,
  set -- $ports
  for port in "$@"; do
    server=$iface:$port
    start editstream/server server=$server redis=$redis ||
      status editstream/server server=$server
  done
end script
