# modules/mathoid/files/mathoid.conf
# vim: set ft=upstart:

# Upstart job configuration for Mathoid
# This file is managed by Puppet

description "TeX to MathML conversion"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid "mathoid"
setgid "mathoid"

env NODE_PATH="/srv/deployment/math/Math/mathoid/node_modules"

chdir "/srv/deployment/math/Math/mathoid"
exec node mathoid

respawn
