# modules/gitblit/files/gitblit.conf
# gitblit - Web interface for Git
#
# Gitblit is an open-source, pure Java stack for managing, viewing, and
# serving Git repositories.
#
description "Web interface for Git"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 4096 4096

setuid gitblit
setgid gitblit

chdir /var/lib/gitblit

exec /usr/bin/java -server -Xmx8g -Djava.awt.headless=true \
        -jar gitblit.jar --baseFolder /var/lib/gitblit/data

respawn
respawn limit 10 5
