# modules/codesearch/templates/initscripts/codesearch-frontend.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=codesearch-frontend
After=docker.service
BindsTo=docker.service

[Service]
User=codesearch
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker kill codesearch-frontend
ExecStartPre=-/usr/bin/docker rm -f codesearch-frontend
# Connect to the local Hound container directly via the docker host IP.
# Avoid tripping CDN rate limits since all conns would come from the same IP.
# https://phabricator.wikimedia.org/T361899
ExecStart=/usr/bin/docker run -p 3003:80 --name codesearch-frontend \
    -e "CODESEARCH_HOUND_BASE=http://172.17.0.1:3002" \
    codesearch-frontend
ExecStop=/usr/bin/docker stop codesearch-frontend
Restart=on-failure

[Install]
WantedBy=multi-user.target
