# modules/varnishkafka/templates/initscripts/varnishkafka.upstart.erb
# varnishkafka - Varnish log collector with Apache Kafka integration

description "Varnish log collector"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid varnishlog
setgid varnishlog

exec /usr/bin/varnishkafka -S "/etc/varnishkafka/<%= @name %>.conf"

respawn
# Seems broken on precise:
# respawn limit unlimited

# vim: set ft=upstart:
