# modules/k8s/templates/initscripts/flannel.upstart.erb
description "flannel overlay network"

start on started network-services

setuid root
setgid root

exec /usr/bin/flanneld \
    --etcd-endpoints=<%= @etcd_endpoints %> \
    --etcd-cafile=/var/lib/puppet/ssl/certs/ca.pem \
    --ip-masq=true

respawn
