# modules/k8s/templates/initscripts/kube-proxy.upstart.erb
description "Kubernetes Proxy"

start on started network-services

setuid root
setgid root

exec /usr/bin/kube-proxy \
    --kubeconfig=/etc/kubernetes/kubeconfig \
    --proxy-mode=iptables \
    --masquerade-all=true

respawn
