# modules/k8s/templates/initscripts/kube-proxy.upstart.erb
description "Kubernetes Proxy"

start on started network-services

setuid root
setgid root

<%-
daemon_args = ['--kubeconfig=/etc/kubernetes/kubeconfig']
if @proxy_mode
	daemon_args.push("--proxy-mode=#{@proxy_mode}")
end
if @masquerade_all
	daemon_args.push("--masquerade-all=true")
end
-%>
exec /usr/bin/kube-proxy \
    <%= daemon_args.join(" \\\n    ") %>

respawn
