# modules/sysctl/files/procps-puppet.conf
# procps-puppet - set puppet-managed sysctls
#
# This task sets Puppet-managed sysctl kernel parameters
# from /etc/sysctl.d/puppet-managed.
description "set sysctls from /etc/sysctl.d/puppet-managed"
start on stopping procps

task

script
    if [ -x /sbin/sysctl ]; then
        for config in /etc/sysctl.d/puppet-managed/*.conf; do
            [ -e "$config" ] || break
            sysctl -e -p "$config"
        done
    fi
end script
