# modules/libraryupgrader/templates/initscripts/libup-web.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=libup web gunicorn daemon
After=network.target

[Service]
PIDFile=/run/gunicorn/pid
User=libup
Group=libup
RuntimeDirectory=gunicorn
WorkingDirectory=/srv/libraryupgrader
ExecStart=/srv/libraryupgrader/venv/bin/gunicorn --pid /run/gunicorn/pid   \
          -w 4 -b 0.0.0.0:3002 libup.web:app
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
RuntimeMaxSec=86400
Restart=on-failure
PrivateTmp=true

[Install]
WantedBy=multi-user.target
