# modules/libraryupgrader/templates/initscripts/libup-ssh-agent.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
[Unit]
Description=SSH key agent

[Service]
Type=simple
User=libup
Group=libup
Environment=SSH_AUTH_SOCK=/tmp/ssh-agent.socket
ExecStart=/usr/bin/ssh-agent -D -a $SSH_AUTH_SOCK

[Install]
WantedBy=multi-user.target
