# modules/redis/templates/initscripts/redis-instance.upstart.erb
description "Advanced key-value store"

start on started network-services

setuid redis
setgid redis

<% unless scope['redis::instance::settings']['daemonize'] == false %>expect fork<% end %>
exec /usr/bin/redis-server /etc/redis/<%= scope['redis::instance::instance_name'] %>.conf

respawn

# vim: set ft=upstart:
