# modules/redis/templates/initscripts/redis-startup.upstart.erb
description "create required directories for redis"

start on mounted MOUNTPOINT=/run TYPE=tmpfs

task

exec install -d -o redis -g redis /var/run/redis
