# modules/mesos/files/marathon.conf
description "Marathon scheduler for Mesos"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

pre-start script
    if [ -f /etc/default/marathon ]; then
        /etc/default/marathon
    fi
    if [ -f /etc/marathon-private ]; then
        /etc/default/marathon-private # For passwords
    fi
end script

exec /usr/bin/marathon
