# modules/ocg_collection/files/ocg-collection.conf
# vim: set ft=upstart:

#####################################################################
### THIS FILE IS MANAGED BY PUPPET
### /modules/ocg_collection/templates/ocg-collection.conf
#####################################################################

# Upstart configuration file for the MediaWiki Collection extension OCG

description "MediaWiki Collection Offline Content Generator"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid ocg-collection
setgid ocg-collection

script
    PATH=$PATH:"<%= @deploy_path %>/node_modules/.bin"
    mediawiki-collection-ocg -c /etc/ocg-collection.js
end script

respawn
