# modules/ocg_collection/files/ocg-collection.conf
# vim: set ft=upstart:

#####################################################################
### THIS FILE IS MANAGED BY PUPPET
### puppet:///ocg_collection/ocg-collection.conf
#####################################################################

# Upstart configuration file for the MediaWiki Collection extension OCG

description "MediaWiki Collection Offline Content Generator"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid ocg-collection
setgid ocg-collection

exec /srv/deployment/ocg_collection/mw-collection-ocg.js -c /etc/ocg-collection.js

respawn
