# modules/openstack2/templates/initscripts/designate-pool-manager.upstart.erb
description "Designate pool manager"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

setuid designate

chdir /var/lib/designate
exec /usr/bin/designate-pool-manager --config-file=/etc/designate/designate.conf --log-file=/var/log/designate/designate-pool-manager.log
respawn
respawn limit 10 30
