# modules/mariadb/files/pt-heartbeat.service
# This file is managed by Puppet

[Unit]
After=mysql.service

[Service]
ExecStart=/usr/bin/pt-heartbeat --update
User=heartbeat

