# modules/xenon/files/xenon-log.conf
# xenon-log - Log Xenon traces to disk.

description "Xenon trace logger"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid xenon
setgid xenon

respawn
respawn limit unlimited

exec /usr/bin/python /usr/local/bin/xenon-log /etc/xenon-log.yaml
