# modules/statsite/files/init/init.conf
description "Start statsite instances"
author "Filippo Giunchedi <filippo@wikimedia.org>"

start on runlevel [2345]

task

script
  /sbin/statsitectl start
end script
