# modules/statsite/files/init/init.conf
description "Start statsite instances"
author "Filippo Giunchedi <filippo@wikimedia.org>"

start on runlevel [2345] or statsite.start

task

script
  for config in /etc/statsite/*.ini; do
    [ -e "$config" ] || break
    name="$(basename $config .ini)"
    start statsite/instance NAME="$name" CONFIG="$config" ||
      status statsite/instance NAME="$name"
  done
end script
