# modules/cgred/templates/initscripts/cgrulesengd.sysvinit.erb
#!/bin/bash
#
# Start/Stop the CGroups Rules Engine Daemon
#
# cgrulesengd   - GGroups Rules Engine Daemon
# description:  This is a daemon for automatically classifying processes \
#       into cgroups based on UID/GID.
#
### BEGIN INIT INFO
# Provides:     cgrulesengd
# Required-Start:   $local_fs $syslog $cgconfig
# Required-Stop:    $local_fs $syslog
# Should-Start:
# Should-Stop:
# Short-Description:    start and stop the cgroups rules engine daemon
# Description:      CGroup Rules Engine is a tool for automatically using \
#                   cgroups to classify processes

SCRIPTNAME=/etc/init.d/$NAME
CGRED_BIN=/usr/sbin/cgrulesengd
CGCONFIG_BIN=/usr/sbin/cgconfigparser
CGCONFIG_DIR=/etc/cgconfig.d/
CGRED_DIR=/etc/cgrules.d/
OPTIONS=""

# Sanity checks
[ -x $CGRED_BIN ] || exit 1
[ -x $CGRED_DIR ] || exit 1
[ -x $CGCONFIG_BIN ] || exit 1
[ -x $CGCONFIG_DIR ] || exit 1

CGRED_DIR_PERM=$(/usr/bin/stat -c %a $CGRED_DIR)
if [ "$CGRED_DIR_PERM" != "555" ]; then
    echo "${CGCONFIG_DIR} permissions are not 500"
    exit 1
fi

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

RETVAL=0

start()
{
    echo $"Starting CGroup Rules Engine Daemon..."
        /bin/cat $CGRED_DIR/*.conf > /etc/cgrules.conf
        $CGCONFIG_BIN -L $CGCONFIG_DIR && \
            start-stop-daemon -S -x $CGRED_BIN $OPTIONS
    RETVAL=$?
}

stop()
{
    echo -n $"Stopping CGroup Rules Engine Daemon..."
        start-stop-daemon -K -x $CGRED_BIN
        echo
    RETVAL=$?
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        start-stop-daemon -T -x $CGRED_BIN
        [ "$?" = 0 ] && echo "OK"
        RETVAL=$?
        ;;
    restart)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart}"
        ;;
esac

exit $RETVAL
