# modules/toollabs/templates/toolschecker.upstart.erb
# This file is managed by Puppet
description "toolchecker UWSGI application"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid <%= @labsproject %>.toolschecker
setgid www-data

# Needs the '=', otherwise doesn't work for chmod-socket. WTF?
# Only 1 worker, since we don't want to run parallel commands at gridengine
# Enforce a timeout of 60s for all our stuff, to prevent tests from clogging other tests
exec /usr/bin/uwsgi --plugin=python \
                    --master \
                    --workers=1 \
                    --harakiri=300 \
                    --harakiri-verbose \
                    --wsgi-file=/usr/local/lib/python2.7/dist-packages/toolschecker.py \
                    --callable=app \
                    --socket=/run/toolschecker/toolschecker.sock \
                    --chmod-socket=664 \
                    --die-on-term

respawn
respawn limit 10 5

# vim: set ft=upstart:
