# modules/toollabs/templates/updatetools.conf.erb
# updatetools: Updates tools and users tables

description "updatetools: Updates tools and users tables"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

console log

setuid <%= @labsproject %>.admin
setgid <%= @labsproject %>.admin

exec /usr/local/bin/updatetools
