# modules/mediawiki/templates/monitoring/mwerrors.conf.erb
# mwerrors -- report MediaWiki exceptions / fatals to StatSD

description "MediaWiki errors StatsD reporter"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid nobody
setgid nogroup

respawn

exec /usr/local/bin/mwerrors           \
    --listen-port <%= @listen_port %>  \
    --statsd-host <%= @statsd_host %>  \
    --statsd-port <%= @statsd_port %>

# vim: set ft=upstart:
