# modules/ircd/files/upstart/ircecho.conf
description "start and stop the irc echo server"
version "1"
author "Chase Pettet"

respawn limit 15 5

start on startup
stop on shutdown

console log

respawn

exec start-stop-daemon --start --make-pidfile --chuid irc --pidfile /var/run/ircecho.pid \
    --exec /usr/local/bin/udpmxircecho.py rc-eqiad argon.wikimedia.org
