# modules/ganglia_new/files/upstart/ganglia-monitor-aggregator.conf
# files/ganglia-monitor-aggregator.conf

description "Ganglia Monitor aggregator"
author "Mark Bergsma <mark@wikimedia.org>"

start on runlevel [2345]

task

script
	for gmonconf in /etc/ganglia/aggregators/*.conf
	do
		gmonid=${gmonconf##*/}
		if status ganglia-monitor-aggregator-instance ID=${gmonid%.conf} && [ /etc/ganglia/aggregators/$gmonid -nt /var/run/gmond-${gmonid%.conf}.pid ]
		then
				restart ganglia-monitor-aggregator-instance ID=${gmonid%.conf} || true
		else
				start ganglia-monitor-aggregator-instance ID=${gmonid%.conf} || true
		fi
	done
end script
