# modules/swift/templates/swiftrepl-mw.service.erb
[Unit]
Description=Ensure mediawiki containers are synchronized across sites
After=network.target

[Service]
User=swiftrepl
ExecStart=/usr/local/bin/swiftrepl-mw repl commons notcommons unsharded global timeline transcoded
