# files/upstart/dynamic-proxy-api.conf
description	"API for the dynamic proxy, aka Invisible Unicorn"

start on runlevel [2345]

respawn
respawn limit 10 5

/usr/bin/dynamic-proxy-api
