# files/upstart/twemproxy.conf
description "twemproxy"

start on (net-device-up
        and runlevel [2345])
stop on runlevel [!2345]

pre-start script
        [ -f /etc/wikimedia-site ] || { stop; exit 0; }
        [ -f "/usr/local/apache/common/wmf-config/twemproxy-$(cat /etc/wikimedia-site)".yaml ] || { stop; exit 0; }
end script

exec /usr/local/bin/nutcracker -m 65536 -c "/usr/local/apache/common/wmf-config/twemproxy-$(cat /etc/wikimedia-site).yaml"

