# files/upstart/varnishhtcpd.conf
# varnishhtcpd - varnish HTCP daemon
#
# Receives multicast HTCP UDP packets and sends corresponding HTTP PURGE requests to Varnish

description "varnish HTCP daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env DEFAULTFILE=/etc/default/varnishhtcpd

script
	DAEMON_OPTS=""
	if [ -f "$DEFAULTFILE" ]; then
	. "$DEFAULTFILE"
	fi
	
	exec /usr/local/bin/varnishhtcpd --name=varnishhtcpd --foreground $DAEMON_OPTS
end script