# modules/wdqs/templates/initscripts/wdqs-blazegraph.upstart.erb
description "WDQS Blazegraph service"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

setuid <%= @username %>
console log
chdir <%= @package_dir %>

exec /bin/bash <%= @package_dir %>/runBlazegraph.sh
respawn
respawn limit 10 30