# modules/openstack/templates/initscripts/designate-mdns.upstart.erb
description "Designate mdns"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

setuid designate

chdir /var/lib/designate
exec /usr/bin/designate-mdns --config-file=/etc/designate/designate.conf --log-file=/var/log/designate/designate-mdns.log
respawn
respawn limit 10 30
