# modules/openstack/templates/initscripts/nova-fullstack.upstart.erb
# nova-fullstack-otest

description "test instance lifecycle (create/setup/delete)"
author "Chase Pettet cpettet@wikimedia.org"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# respawn the job up to 10 times within a 5 second period.
# If the job exceeds these values, it will be stopped and
# marked as failed.
respawn
respawn limit 10 60

script
    export OS_USERNAME="osstackcanary"
    export OS_PASSWORD="<%= @password %>"
    export OS_REGION_NAME="<%= @region %>"
    /usr/local/sbin/nova-fullstack \
                                       --ssh-timeout <%= @ssh_timeout %> \
                                       --puppet-timeout <%= @puppet_timeout %> \
                                       --creation-timeout <%= @creation_timeout %> \
                                       --interval <%= @interval %> \
                                       --max-pool <%= @max_pool %>  \
                                       --project admin-monitoring \
                                       --prepend fullstackd \
                                       --image debian-9.8-stretch \
                                       --flavor m1.small \
                                       --keyfile <%= @keyfile %> \
                                       --statsd statsd.eqiad.wmnet
end script
