# modules/abacist/templates/abacist.conf.erb
# abacist
#
# This is an Upstart job configuration file for Abacist, a redis-backed
# web analytics aggregator. This file was provisioned by Puppet.

description "web analytics aggregator"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid abacist
setgid abacist

respawn
respawn limit unlimited

exec /usr/bin/python /srv/deployment/abacist/abacist/abacist.py <%= @eventlogging_publisher %> --redis-server=<%= @redis_server %>
