# modules/webperf/files/asset-check.conf
# vim: set ft=upstart:

# Upstart job configuration for asset-check.py
# This file is managed by Puppet
description "Report static asset stats to Ganglia"

start on (local-filesystems and net-device-up IFACE!=lo)

setuid webperf
setgid webperf

chdir /srv/webperf
exec python asset-check.py

respawn
