# modules/ganglia/files/systemd/ganglia-monitor-aggregator.service
[Unit]
Description=Ganglia aggregator

[Service]
TimeoutStartSec=0
ExecStart=/usr/sbin/gmond

[Install]
WantedBy=multi-user.target

