# modules/phabricator/templates/initscripts/aphlict.upstart.erb
# aphlict - Phabricator notification server
#
# Routes real-time notifications to websocket clients

description "Phabricator notification websocket server"

setuid <%= scope.lookupvar('phabricator::aphlict::user') %>
setgid <%= scope.lookupvar('phabricator::aphlict::group') %>

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

expect stop

exec <%= scope.lookupvar('phabricator::aphlict::aphlict_cmd') %>
