# modules/apertium/templates/initscripts/apertium-apy.upstart.erb
#####################################################################
### THIS FILE IS MANAGED BY PUPPET
#####################################################################

description "Apertium APY service"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid "<%= @uid %>"
setgid "<%= @gid %>"

# Apertium seems to leak fds, or the current default limit is not enough
limit nofile 8096 8096

env LC_ALL=en_US.utf-8

respawn
# Warning: daemon mode for apertium is not what you might think it is. From the list
# of things to do for correct daemon behaviour in
# https://www.python.org/dev/peps/pep-3143/#correct-daemon-behaviour
# it only seems to do the first. So we abuse it a bit to get standard behaviour
# and better logging
exec /usr/bin/python3 /usr/share/apertium-apy/servlet.py -j<%= @num_of_processes %> -m<%= @max_idle_seconds %> --daemon --log-path <%= @log_dir %> <%= @mode_path %>

# wait 60 seconds for a graceful restart before killing the master
kill timeout 60
