# modules/mw_rc_irc/files/upstart/ircd.conf
description "start and stop the ircd server"
version "1"
author "Daniel Zahn"

setuid irc
setgid irc

start on startup
stop on shutdown

console log

exec /usr/bin/ircd -foreground
