# modules/samplicator/templates/samplicator.service.erb
<%#- SPDX-License-Identifier: Apache-2.0 -%>
# Managed by Puppet
[Unit]
Description=Samplicator
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/samplicate -b <%= @recvbuf %> -S -p <%= @port %> -d 0 -4 -f <%= @targets.join(' ') %>
Restart=on-failure
RestartSec=0
SyslogIdentifier=samplicator

[Install]
WantedBy=multi-user.target
