# modules/mediawiki/files/hhvm/ensure_jemalloc_deactivated.cron
# This file is managed by Puppet
*/30 * * * * apache /usr/local/bin/deactivate-jemalloc >/dev/null 2>&1
