# modules/visualdiff/files/parsoid-vd-client.upstart.conf
# vim: set ft=upstart:

# Upstart job configuration for visual diffing client

description "Visual diffing client"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

setuid "testreduce"
setgid "testreduce"

chdir "/usr/lib/testreduce/client"
exec node rtclient-cluster.js -c 4 /srv/visualdiff/testreduce/parsoid-vd-client.config.js

respawn
