# modules/openstack/files/ocata/designate/designate-producer.conf
# This init scripts are not present in the Newton debian packages.
#  Probably we won't need this on Ocata.

description "Designate producer"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

setuid designate

chdir /var/lib/designate
exec /usr/bin/designate-producer --config-file=/etc/designate/designate.conf --log-file=/var/log/designate/designate-producer.log
respawn
respawn limit 10 30
